# scripts/ubuntu/scrabble-ninja.conf
description "scrabble-ninja upstart job"
author "Rob Polak"
start on runlevel [2345]
exec su root -c /var/node/scrabble-ninja/run.sh
